# crontab
# Run every day at 19:00
0 19 * * * /app/bluesnews-full-rss --mode=fetch >> /var/log/cron.log 2>&1
